# supervisor-ipython.conf
[program:ipython]
command=/usr/local/bin/ipython notebook
user=vagrant

