# supervisord-vagrant.conf
[program:ipython]
command=/usr/local/bin/ipython notebook --ip=* --notebook-dir=/notebooks --pylab inline
user=vagrant

