# vagrant/supervisor-ipython.conf
[program:ipython]
command=/usr/local/bin/ipython notebook --ip=0.0.0.0
user=vagrant
directory=/data

