# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:meqserver]
command=/usr/bin/meqserver

[program:ipython]
command=ipython notebook --ip=* \
 --MappingKernelManager.time_to_dead=10 \
 --MappingKernelManager.first_beat=3 --notebook-dir=/notebooks