# meqtrees-server/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:meqserver]
command=/usr/bin/meqserver
